# LunaReady.upstart
# -*- mode: shell-script; -*-

start on LunaSysMgr-ready

script
         setcpushares-task -p $(pidof WebAppMgr) -t WebAppMgr || true
end script
